# mtd-agent.service
[Unit]
Description=MTD Edge Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/mtd-agent/agent
ExecStart=/opt/mtd-agent/venv/bin/python main.py
Restart=always
RestartSec=10
Environment=MTD_CONFIG=/opt/mtd-agent/config.json
Environment=MTD_DB=/opt/mtd-agent/cache.db
Environment=MTD_API_URL=https://api.mijnthuisdata.nl

[Install]
WantedBy=multi-user.target
