# systemd/mtd-core.service
[Unit]
Description=MTD Core (heartbeat, WebSocket, OTA-updates, statuspagina)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/mtd-agent/agent
ExecStart=/opt/mtd-agent/venv/bin/python core.py
Restart=always
RestartSec=5
Environment=MTD_CONFIG=/opt/mtd-agent/config.json
Environment=MTD_STATE=/opt/mtd-agent/state.json
Environment=MTD_SIGNALS=/opt/mtd-agent/signals.json
Environment=MTD_PLUGIN_DIR=/opt/mtd-agent/agent/integrations
# Optioneel per-device bestand om MTD_API_URL/MTD_WS_URL te overrulen (buiten
# de git working tree, overleeft dus OTA-updates). Ontbreekt het, dan geldt de
# fallback-default uit agent/api.py en agent/websocket_client.py.
EnvironmentFile=-/etc/mtd-agent/env

[Install]
WantedBy=multi-user.target
