# systemd/mtd-portal.service
[Unit]
Description=MTD Captive Portal
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/mtd-agent/captive_portal
ExecStart=/opt/mtd-agent/venv/bin/python portal.py
Restart=on-failure
Environment=MTD_CONFIG=/opt/mtd-agent/config.json
# Optioneel per-device bestand om MTD_API_URL te overrulen (zie mtd-core.service).
EnvironmentFile=-/etc/mtd-agent/env

[Install]
WantedBy=multi-user.target
