# systemd/mtd-agent.service
[Unit]
Description=MTD Agent (edge agent: plugins, WS-sync, lokale webserver)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/mtd-agent
ExecStart=/opt/mtd-agent/venv/bin/python main.py
Restart=always
RestartSec=5
Environment=DB_PATH=/data/agent.db
Environment=PLUGIN_DIR=/data/plugins
# Optioneel per-device bestand om PLATFORM_WS_URL/PLATFORM_API_URL/AGENT_KEY te
# overrulen (buiten de git working tree, overleeft dus OTA-updates).
EnvironmentFile=-/etc/mtd-agent/env

[Install]
WantedBy=multi-user.target
