# systemd/mtd-portal.service
[Unit]
Description=MTD Captive Portal (onboarding)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/mtd-agent
ExecStart=/opt/mtd-agent/venv/bin/python -m onboarding.portal
Restart=on-failure
Environment=DB_PATH=/data/agent.db
EnvironmentFile=-/etc/mtd-agent/env

[Install]
WantedBy=multi-user.target
